# scripts/.backup/scheduled_tasks/QM_PublicSnapshot_Export_Hourly.2026-05-12_pre_L2.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.3" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Exports public-data snapshot JSON hourly and publishes if changed.</Description>
    <URI>\QM_PublicSnapshot_Export_Hourly</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <IdleSettings>
      <Duration>PT10M</Duration>
      <WaitTimeout>PT1H</WaitTimeout>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Triggers>
    <TimeTrigger>
      <StartBoundary>2026-05-08T00:07:00+02:00</StartBoundary>
      <Repetition>
        <Interval>PT1H</Interval>
        <Duration>P3650D</Duration>
        <StopAtDurationEnd>true</StopAtDurationEnd>
      </Repetition>
    </TimeTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -WindowStyle Hidden -File "C:\QM\paperclip\data\instances\default\projects\03d4dcc8-4cea-4133-9f68-90c0d99628fb\ac8daa03-00ae-49fd-bd4a-f1283a075f83\_default\scripts\export_public_snapshot.ps1" -NoGit</Arguments>
      <WorkingDirectory>C:\QM\paperclip\data\instances\default\projects\03d4dcc8-4cea-4133-9f68-90c0d99628fb\ac8daa03-00ae-49fd-bd4a-f1283a075f83\_default</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
